# files/etc/NetworkManager/dispatcher.d/00firewall
#!/bin/sh

set -e

case "$2" in
    up)
        [ -x /usr/sbin/ferm ] || exit 2
        /bin/systemctl reload ferm
    ;;        
    *)
    exit 0
    ;;
esac

